# urbackupclientbackend-redhat.service
[Unit]
Description=UrBackup Client backend
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/urbackupclientbackend --config /etc/sysconfig/urbackupclientbackend --no-consoletime
User=root

[Install]
WantedBy=multi-user.target

